# dev/config/systemd/user/backup_weekly.service
[Unit]
Description=Weekly OS Backup

[Service]
Type=oneshot
ExecStart=%h/.backup_weekly.sh

