# dev/config/systemd/user/backup_weekly.timer
[Unit]
Description=Run Weekly OS Backup

[Timer]
OnCalendar=Mon
Persistent=true

[Install]
WantedBy=timers.target

