# dev/config/systemd/user/backup_daily.service
[Unit]
Description=Daily OS Backup

[Service]
Type=oneshot
ExecStart=%h/.backup_daily.sh

